# deploy/macos/com.askdao.imessage-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identity -->
    <key>Label</key>
    <string>com.askdao.imessage-bridge</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <!-- Command: npx tsx src/index.ts -->
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__/npx</string>
        <string>tsx</string>
        <string>src/index.ts</string>
    </array>

    <!-- Environment variables loaded from .env by dotenv, but PATH is needed -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__NODE_PATH__:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Auto-start on load (boot/login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Auto-restart on crash -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Throttle restart: wait 10s before restarting -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/logs/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/logs/bridge.error.log</string>
</dict>
</plist>
